# srv/salt/base/supervisor/files/supervisor.conf
[program:hello]
command = gunicorn hello:app
directory = /usr/share/nginx/www/flask/
user = root
